# etc/systemd/system/chatgpt.service
[Unit]
Description=ChatGPT Telegram Bot 2.0
After=network-online.target
Wants=network-online.target

[Service]
User=chatgpt
Group=chatgpt
UMask=0027

WorkingDirectory=/home/chatgpt
EnvironmentFile=/etc/chatgpt.env

ExecStart=/home/chatgpt/venv/bin/python /home/chatgpt/bot.py

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/chatgpt/state
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
PrivateDevices=yes
PrivateTmp=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_INET

Restart=on-failure
RestartSec=10

TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
